# docker/supervisord.conf
[supervisord]
# Ejecuta en foreground para que Docker pueda supervisar el proceso principal.
nodaemon=true
# Archivos de log y PID de supervisord.
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
# Servidor HTTP de Laravel que recibe las peticiones del proxy inverso.
command=php artisan serve --host=0.0.0.0 --port=8000
autostart=true
autorestart=true
# Logs stdout/stderr del servidor Laravel.
stdout_logfile=/var/log/laravel-http.log
stderr_logfile=/var/log/laravel-http-error.log

[program:nginx]
# Nginx en foreground (requerido dentro de contenedor).
command=nginx -g "daemon off;"
autostart=true
autorestart=true
# Logs stdout/stderr del proceso nginx.
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log

# Configuracion anterior para arrancar PHP-FPM en lugar del servidor HTTP de Laravel.
# [program:php-fpm]
# command=php-fpm
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/php-fpm.log
# stderr_logfile=/var/log/php-fpm-error.log