# docker/supervisord.conf
[supervisord]
# Ejecuta en foreground para que Docker pueda supervisar el proceso principal.
nodaemon=true
# Archivos de log y PID de supervisord.
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
# Proceso PHP-FPM que ejecuta scripts PHP.
command=php-fpm
autostart=true
autorestart=true
# Logs stdout/stderr del proceso php-fpm.
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-error.log

[program:nginx]
# Nginx en foreground (requerido dentro de contenedor).
command=nginx -g "daemon off;"
autostart=true
autorestart=true
# Logs stdout/stderr del proceso nginx.
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log